# enterprise_rag_system/deployment/systemd/enterprise-rag-backend.service
[Unit]
Description=Enterprise RAG FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=enterprise-rag
Group=enterprise-rag
WorkingDirectory=/opt/enterprise-rag/app/backend
EnvironmentFile=/etc/enterprise-rag/production.env
ExecStart=/opt/enterprise-rag/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
